# Prometheus/Ubuntu-20.04-LTS/cron.d/my_fake_data

*/5 * * * * node_exporter /prometheus/code/prometheus_random.py
